# workshop-october-2024/demo/services/demo-quadlet-monitor.container
[Unit]
Description=Tire Monitor

[Container]
Image=localhost/tire-monitor
StopSignal=SIGINT

Network=host
PodmanArgs=--pull=never

[Service]
Restart=no
