# workshop-october-2024/demo/services/demo-tire-sensor.service
[Unit]
Description=Tire Sensor Service

[Service]
Type=simple
ExecStart=python3 /var/demo/tire-sensor.py
